# my.cron
@reboot /usr/bin/forever start /lightsail-client/app.js